# tproxy-route.service
[Unit]
Description=TPROXY NAT Gateway – Policy Routing Setup
Documentation=https://wiki.nftables.org/
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes

# 参数: fwmark 路由表 网卡 子网
#       必须与 nftables-tproxy.conf 中的定义保持一致
ExecStart=/usr/local/bin/setup-tproxy-route.sh 1 100 wlo1 192.168.2.0/24

ExecStop=/bin/bash -c '\
  ip -4 rule  del fwmark 1 table 100 2>/dev/null; \
  ip -4 route del local 0.0.0.0/0 dev lo table 100 2>/dev/null; \
  sysctl -w net.ipv4.ip_forward=0 >/dev/null; \
  echo "TPROXY NAT gateway routing removed."'

[Install]
WantedBy=multi-user.target
